# host/com.myrunner.pool.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN"
  "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<!--
launchd service: keeps the runner pool alive across reboots.

Install:
  cp host/com.myrunner.pool.plist ~/Library/LaunchAgents/
  launchctl load ~/Library/LaunchAgents/com.myrunner.pool.plist

Uninstall:
  launchctl unload ~/Library/LaunchAgents/com.myrunner.pool.plist
  rm ~/Library/LaunchAgents/com.myrunner.pool.plist

Update REPO_PATH below to the absolute path of this repository.

launchd has no keyring access: configure GitHub App credentials in .env
(GITHUB_APP_ID / GITHUB_APP_INSTALLATION_ID / GITHUB_APP_PRIVATE_KEY_PATH).
Do not rely on `gh auth token` for this service.
-->
<plist version="1.0">
<dict>
  <key>Label</key>
  <string>com.myrunner.pool</string>

  <key>ProgramArguments</key>
  <array>
    <string>/bin/bash</string>
    <string>REPO_PATH/scripts/pool.sh</string>
  </array>

  <key>WorkingDirectory</key>
  <string>REPO_PATH</string>

  <key>RunAtLoad</key>
  <true/>

  <key>KeepAlive</key>
  <true/>

  <key>StandardOutPath</key>
  <string>REPO_PATH/pool.log</string>

  <key>StandardErrorPath</key>
  <string>REPO_PATH/pool.log</string>

  <key>EnvironmentVariables</key>
  <dict>
    <key>PATH</key>
    <string>/opt/homebrew/bin:/usr/local/bin:/usr/bin:/bin</string>
  </dict>
</dict>
</plist>
